# packer/files/docker-veth-fix.service
[Unit]
Description=Fix Docker bridge veth attachment (VirtualBox workaround)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/docker-veth-fix.sh
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
